# systemd/claude-agent.service
[Unit]
Description=Claude Automation Agent
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory={{REPO_DIR}}/agent
ExecStart={{VENV_DIR}}/bin/python claude-agent.py
EnvironmentFile={{REPO_DIR}}/.env
Environment=DATA_DIR={{DATA_DIR}}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=claude-agent

# Limits
TimeoutStartSec=300
MemoryMax=2G
